# .config/autostart/randr.desktop
[Desktop Entry]
Type=Application
Name=Randr Setup
Exec=~/bin/Xsetup
OnlyShowIn=Awesome;
X-GNOME-Autostart-enabled=true
NoDisplay=true
